# bb.conf
[program:bb]
command=python3 /home/curator2/bb/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/bb/test.err.log
stdout_logfile=/home/curator2/bb/test.out.log
